# scripts/com.dex.bp-citations.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dex.bp-citations</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dex/YDEX/INPORTANT WORK/BACKLINK/backlink-publisher-69/backlink-publisher/scripts/run-citations-periodic.sh</string>
    </array>
    <!-- Daily 06:00 (after recheck 04:30, selector-drift 05:00).
         Conservative quota defaults: --max-pairs 5 / --cost-cap 10 per run.
         Adjust BP_CITATIONS_MAX_PAIRS / BP_CITATIONS_COST_CAP after confirming
         Perplexity v1 API daily quota. See docs/operations/probe-citations-runbook.md. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/dex/YDEX/INPORTANT WORK/BACKLINK/backlink-publisher-69/backlink-publisher</string>
    <key>StandardOutPath</key>
    <string>/Users/dex/YDEX/INPORTANT WORK/BACKLINK/backlink-publisher-69/backlink-publisher/logs/citations-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dex/YDEX/INPORTANT WORK/BACKLINK/backlink-publisher-69/backlink-publisher/logs/citations-launchd.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <!-- Adjust these two after confirming Perplexity v1 daily quota.
             Default: max 5 pairs / 10 probes per run (safe conservative start). -->
        <key>BP_CITATIONS_MAX_PAIRS</key>
        <string>5</string>
        <key>BP_CITATIONS_COST_CAP</key>
        <string>10</string>
    </dict>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
