# scripts/com.dex.bp-keepalive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dex.bp-keepalive</string>

    <!-- keepalive-run chains recheck → gap → publish → reverify automatically.
         The cycle-level lock (.keepalive-run.lock) prevents concurrent runs if
         a cycle exceeds 24h; each step has its own internal budget (_BATCH_BUDGET_S).

         Scheduling note: keepalive-run fires at 06:30, after full-pipeline (04:00).
         collect-signals uses StartInterval (interval-based, not calendar-pinned) so
         no strict ordering between keepalive and collect-signals is guaranteed.

         Install:
           cp scripts/com.dex.bp-keepalive.plist ~/Library/LaunchAgents/
           launchctl load ~/Library/LaunchAgents/com.dex.bp-keepalive.plist

         Uninstall:
           launchctl unload ~/Library/LaunchAgents/com.dex.bp-keepalive.plist
           rm ~/Library/LaunchAgents/com.dex.bp-keepalive.plist

         Logs: tail -f logs/keepalive-launchd.log
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/dex/YDEX/INPORTANT WORK/外链/backlink-publisher-69/backlink-publisher/.venv/bin/python</string>
        <string>-m</string>
        <string>backlink_publisher.cli.keepalive_run</string>
    </array>

    <!-- Daily at 06:30 local time.
         Full-pipeline (com.dex.bp-full-pipeline.plist) runs at 04:00 with a
         90-minute batch budget — keepalive-run fires after it should be done.
         Both target sticky platforms but use different gap logic (keepalive-gap
         targets dead links; plan-gap targets fresh links); they are complementary. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/dex/YDEX/INPORTANT WORK/外链/backlink-publisher-69/backlink-publisher</string>

    <key>StandardOutPath</key>
    <string>/Users/dex/YDEX/INPORTANT WORK/外链/backlink-publisher-69/backlink-publisher/logs/keepalive-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dex/YDEX/INPORTANT WORK/外链/backlink-publisher-69/backlink-publisher/logs/keepalive-launchd.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>PYTHONHASHSEED</key>
        <string>0</string>
        <key>PYTHONPATH</key>
        <string>/Users/dex/YDEX/INPORTANT WORK/外链/backlink-publisher-69/backlink-publisher/src</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
